# systemd/babymonitor-kiosk.service
[Unit]
Description=BabyMonitor Kiosk (native GTK+GStreamer)
After=babymonitor-monitor.service graphical.target
Requires=graphical.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=KIOSK_MODE=1
# Wait until the monitor FSM has written the camera URL
ExecStartPre=/bin/bash -c 'until [ -f /etc/babymonitor/kiosk_url ]; do sleep 1; done'
ExecStart=/usr/bin/python3 /opt/babymonitor/babymonitor/monitor_app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
