# systemd/babymonitor-kiosk.service
[Unit]
Description=BabyMonitor Kiosk (Chromium)
After=babymonitor-monitor.service graphical.target
Requires=graphical.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
# Wait until the monitor FSM has written the camera URL
ExecStartPre=/bin/bash -c 'until [ -f /etc/babymonitor/kiosk_url ]; do sleep 1; done'
ExecStart=/bin/bash -c '\
    URL=$(cat /etc/babymonitor/kiosk_url); \
    chromium-browser \
        --kiosk \
        --noerrdialogs \
        --disable-infobars \
        --disable-restore-session-state \
        --disable-session-crashed-bubble \
        --disable-translate \
        --disable-features=TranslateUI \
        --touch-events=enabled \
        --enable-features=OverlayScrollbar \
        --check-for-update-interval=31536000 \
        "$URL"'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
